# com.gitorizer.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gitorizer.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tatu.lahtela/.local/bin/gitorizer</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tatu.lahtela/.local/share/gitorizer/gitorizer.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tatu.lahtela/.local/share/gitorizer/gitorizer.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/tatu.lahtela</string>
        <key>PATH</key>
        <string>/Users/tatu.lahtela/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/tatu.lahtela</string>
</dict>
</plist>
